# tres-robo.service
[Unit]
Description=Tres-Robo voice assistant
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Tres-Robo
ExecStart=%h/Tres-Robo/.venv/bin/python main_gemini.py
Restart=always
RestartSec=5
EnvironmentFile=%h/Tres-Robo/.env
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tres-robo

[Install]
WantedBy=default.target
